# launchd/com.snapname.watcher.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.snapname.watcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/ABSOLUTE/PATH/TO/snapname-repo/.venv/bin/python</string>
		<string>-m</string>
		<string>snapname</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/ABSOLUTE/PATH/TO/snapname-repo</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/ABSOLUTE/PATH/TO/snapname-repo/.snapname-launchd.out.log</string>
	<key>StandardErrorPath</key>
	<string>/ABSOLUTE/PATH/TO/snapname-repo/.snapname-launchd.err.log</string>
</dict>
</plist>
